# ansible/roles/web-server/templates/dummy.service.j2
[Unit]
Description= simple webservice to creat network traffic

[Service]
WorkingDirectory=/opt
ExecStart=/usr/bin/python3 dummy-web.server.py 3000
Restart=on-failure
RestartSec=5
