# ansible/roles/install_mongodb/templates/mongodb.service.j2
[Unit]
Description=MongoDB Database Server

[Service]
ExecStart=/usr/bin/mongod --config /etc/mongod.conf
Restart=on-failure
User=mongodb
Group=mongodb
RestartSec=5
